# configs/systemd/openclaw-healthcheck.service
[Unit]
Description=OpenClaw Gateway Healthcheck
After=openclaw-gateway.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/openclaw-healthcheck.sh
