# configs/systemd/openclaw-healthcheck.timer
[Unit]
Description=Run OpenClaw healthcheck every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
Unit=openclaw-healthcheck.service
Persistent=true

[Install]
WantedBy=timers.target
